# contrib/systemd/antyswirusd.service
[Unit]
Description=antyswirus daemon
Documentation=https://github.com/kaykoe/antyswirus/blob/main/README.md
After=local-fs.target
ConditionPathExists=/var/lib/antyswirus

[Service]
Type=forking
PIDFile=/run/antyswirus/antyswirusd.pid
ExecStart=/usr/bin/antyswirusd start
ExecStop=/usr/bin/antyswirusd stop
ExecReload=/bin/sh -c '/usr/bin/antyswirusd stop && /usr/bin/antyswirusd start'

RuntimeDirectory=antyswirus
StateDirectory=antyswirus
LogsDirectory=antyswirus

# Desktop notifications (adjust DISPLAY and UID to match your session)
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
